# supervisord.conf
[supervisord]
logfile=%(here)s/run/supervisor/log
pidfile=%(here)s/run/supervisor/pid
childlogdir=%(here)s/run/supervisor/
nodaemon=true
directory=%(here)s
identifier=cockatoo

[unix_http_server]
file = %(here)s/run/supervisor/sock
chmod = 0770

[supervisorctl]
serverurl=unix://%(here)s/run/supervisor/sock
prompt=cockatoo

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:worker]
command=make run-cuckoo-worker
priority=900
startsecs=60
stopwaitsecs=400
redirect_stderr=true
autorestart=true
startretries=1000

[program:api]
command=make run-cuckoo-dist-api
priority=700
startsecs=30
redirect_stderr=true
autorestart=true
startretries=1000

[program:db]
command=make run-cuckoo-dist-db
priority=600
startsecs=30
redirect_stderr=true
autorestart=true

[group:dist]
programs=api,db
priority=500

[group:cuckoo]
programs=worker
priority=800